# config/files/usr/etc/profile.d/ssh_askpass.sh
export SSH_ASKPASS=/usr/bin/ksshaskpass
export SSH_ASKPASS_REQUIRE=prefer
export PINENTRY_KDE_USE_WALLET=true
